# launchd/com.taiwan-stock-intelligence.weekly-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.taiwan-stock-intelligence.weekly-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/louise87121/anaconda3/bin/python</string>
    <string>/Users/louise87121/Documents/test/Stock Analysis Top50/Taiwan Stock Intelligence/scripts/weekly_refresh.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/louise87121/Documents/test/Stock Analysis Top50/Taiwan Stock Intelligence</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/louise87121/Documents/test/Stock Analysis Top50/Taiwan Stock Intelligence/logs/weekly_refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/louise87121/Documents/test/Stock Analysis Top50/Taiwan Stock Intelligence/logs/weekly_refresh.err</string>
</dict>
</plist>
